# systemd/prd-backend.service
[Unit]
Description=PRD Analyst Backend (FastAPI)
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/path/to/project/backend
Environment="PATH=/usr/bin:/usr/local/bin"
EnvironmentFile=/path/to/project/.env
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=append:/var/log/prd-analyst/backend.log
StandardError=append:/var/log/prd-analyst/backend-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
